# conf.d/supervisor/flasker.conf
[program:application]
directory=/var/www/45.56.119.8
command=gunicorn3 application:application
autostart=true
autorestart=true
stopasgroup=true
killasgroup=true
